# mqtt-node.service
[Unit]
Description=MQTT controlled IO client
After=network.target

[Service]
User=Gerald
Type=simple
WorkingDirectory=/home/Gerald/python-mqtt-node-rpi/
ExecStart=/usr/bin/python -mmqtt_node
StandardError=journal
Restart=always

[Install]
WantedBy=multi-user.target
